# wndr/luci/themes/waio/root/etc/uci-defaults/luci-theme-waio
#!/bin/sh
uci batch <<-EOF
	set luci.themes.WAiO=/luci-static/waio
        commit luci
EOF
	
